# Ansible/roles/clamav/templates/clamd-mediashuttle.timer.j2
[Unit]
Description=ClamAV scan of {{ clamav_directory }}

[Timer]
OnCalendar=OnCalendar=*-*-* 03:00:00 America/Los_Angeles
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
