# Ansible/roles/clamav/templates/clamd@mediashuttle.timer.j2
[Unit]
Description=ClamAV scan of {{ clamav_directory }}

[Timer]
OnCalendar=daily
AccuracySec=1h
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
